# frostbyte.service
[Unit]
Description=FrostByte — auto-suspend inactive apps
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/frostbyte run
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
